# etc/systemd/user/footserver.service
[Unit]
Description=foot server service
Documentation=man:foot

[Service]
ExecStart=/usr/bin/foot --server
