# etc/systemd/user/kanshi.service
[Unit]
Description=Kanshi service
Documentation=man:kanshi

[Service]
ExecStart=/usr/bin/kanshi
