# etc/systemd/user/nm-applet.service
[Unit]
Description=NetworkManager applet service
Documentation=man:nm-applet

[Service]
ExecStart=/usr/bin/nm-applet
