# etc/systemd/user/waybar.service
[Unit]
Description=Waybar service
Documentation=man:waybar

[Service]
ExecStart=/usr/bin/waybar
