# etc/systemd/user/wayland.target
[Unit]
Description=Wayland
RefuseManualStart=no
StopWhenUnneeded=no
